# deploy/systemd/rexfinhub-bloomberg.timer
[Unit]
Description=REX FinHub Bloomberg Pull Timer (every 15 min, 4-6 PM ET weekdays)

[Timer]
# VPS timezone set to America/New_York — DST handled automatically
# Single pull after Ryu finishes Bloomberg updates (~5 PM)
OnCalendar=Mon..Fri 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
