# deploy/systemd/rexfinhub-bloomberg.timer
[Unit]
Description=REX FinHub Bloomberg Pull Timer (once weekdays, 17:15 ET)

[Timer]
# VPS timezone set to America/New_York — DST handled automatically.
# 17:15 gives the Bloomberg xlsm file ~13 minutes past its typical 17:02
# save time to fully populate. Reading at 17:00 was racing the file saver
# and picking up NaN values for tickers like FNGU, causing a silent
# fallback to yesterday's AUM (seen 2026-04-14, $274M under-report).
OnCalendar=Mon..Fri 17:15:00
Persistent=true

[Install]
WantedBy=timers.target
