# deploy/systemd/rexfinhub-daily.timer
[Unit]
Description=REX FinHub Daily Pipeline Timer (weekdays 17:45 ET)

[Timer]
# VPS timezone set to America/New_York — DST handled automatically.
# Bloomberg pull runs at 17:15; 17:45 gives the market-sync + prebake +
# upload chain 30 min of runway before the internal 17:30 email send gate.
# (Note: the email send inside run_daily.py still waits until its own
# target time, so the 17:45 start is the pipeline kickoff, not the send.)
OnCalendar=Mon..Fri 17:45:00
Persistent=true

[Install]
WantedBy=timers.target
